# deploy/systemd/sda-orchestrator-worker@.service
[Unit]
Description=Meraki Dynamic SDA Orchestrator Apply Worker %i
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/sda-orchestrator/current
EnvironmentFile=%h/.config/sda-orchestrator/worker.env
ExecStart=%h/sda-orchestrator/current/.venv/bin/python -m orchestrator.worker_runtime --run-id %i
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=%h/.local/share/sda-orchestrator

[Install]
WantedBy=default.target
